# grafana/grafana.service
[Unit]
Description=Grafana Agent

[Service]
RemainAfterExit=yes
ExecStart=/usr/bin/grafana server

[Install]
WantedBy=multi-user.target